# greenhouse.service
[Unit]
Description=Victory Garden Greenhouse Controller
After=network-online.target mosquitto.service
Wants=network-online.target mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/victory_garden/python_tools
ExecStart=/home/pi/victory_garden/python_tools/.venv/bin/python -m main
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=victory-garden

[Install]
WantedBy=multi-user.target
